# deploy/researchradar.service
[Unit]
Description=ResearchRadar local AI research source agent
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=batchcom
Group=batchcom
WorkingDirectory=/home/dataset-local/ResearchRadar
Environment=PIP_CACHE_DIR=/home/dataset-local/.cache/pip
Environment=HOME=/home/batchcom
Environment=RESEARCHRADAR_HOST=0.0.0.0
Environment=RESEARCHRADAR_PORT=8766
Environment=NO_PROXY=localhost,127.0.0.1,::1,0.0.0.0,127.0.0.0/8
Environment=no_proxy=localhost,127.0.0.1,::1,0.0.0.0,127.0.0.0/8
EnvironmentFile=-/home/dataset-local/ResearchRadar/.env
ExecStart=/usr/bin/env RESEARCHRADAR_HOST=0.0.0.0 RESEARCHRADAR_PORT=8766 NO_PROXY=localhost,127.0.0.1,::1,0.0.0.0,127.0.0.0/8 no_proxy=localhost,127.0.0.1,::1,0.0.0.0,127.0.0.0/8 /home/dataset-local/ResearchRadar/scripts/start.sh
Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target
